# com.plex.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.plex.backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Volumes/WD_BLACK_SN770_1TB/code_claude/PlexBackups/plex_backup.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Volumes/media/backup/PlexBackups/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Volumes/media/backup/PlexBackups/launchd_error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
